# systemd/ipo-check.service
[Unit]
Description=IPO check (run once at 11:11 Nepal time)
After=network.target

[Service]
Type=oneshot
User=IPO_USER
WorkingDirectory=IPO_PROJECT_DIR
Environment="HOME=IPO_HOME"
Environment="PATH=IPO_PROJECT_DIR/venv/bin:/usr/local/bin:/usr/bin:/bin"
ExecStart=/usr/bin/env python3 IPO_PROJECT_DIR/src/scheduler/run_once.py
StandardOutput=journal
StandardError=journal
